# templates/nzbget.service
[Unit]
Description=nzbget
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/nzbget -D
ExecStop=/usr/local/bin/nzbget -Q
Restart=on-failure
User={{ nzbget_user }}

[Install]
WantedBy=multi-user.target
